# boot/jiu-pluck.service
[Unit]
Description=Jiu-Pluck Application (Frontend + Backend)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=PROJECT_ROOT
Environment="PATH=PROJECT_ROOT/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=PROJECT_ROOT/boot/start_production.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

